# deploy/systemd/newscast-web-compose.service
[Unit]
Description=Newscast Navigator web stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/newscast-web
ExecStart=/usr/bin/docker compose --env-file /opt/newscast-web/deploy/env/web-prod.env -f /opt/newscast-web/deploy/docker/docker-compose.web-prod.yml up -d --build
ExecStop=/usr/bin/docker compose --env-file /opt/newscast-web/deploy/env/web-prod.env -f /opt/newscast-web/deploy/docker/docker-compose.web-prod.yml down
ExecReload=/usr/bin/docker compose --env-file /opt/newscast-web/deploy/env/web-prod.env -f /opt/newscast-web/deploy/docker/docker-compose.web-prod.yml up -d --build
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
