# deploy/systemd/newscast-web-compose.service
[Unit]
Description=Newscast Navigator web stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/newscast-web/newscast-web.env
ExecStart=/usr/bin/bash -lc 'cd "${PROJECT_ROOT}" && docker compose --env-file "${COMPOSE_ENV_FILE}" -f "${COMPOSE_FILE}" up -d --build'
ExecStop=/usr/bin/bash -lc 'cd "${PROJECT_ROOT}" && docker compose --env-file "${COMPOSE_ENV_FILE}" -f "${COMPOSE_FILE}" down'
ExecReload=/usr/bin/bash -lc 'cd "${PROJECT_ROOT}" && docker compose --env-file "${COMPOSE_ENV_FILE}" -f "${COMPOSE_FILE}" up -d --build'
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
